# shufflewall/shufflewall.timer
# shufflewall.timer

[Unit]
Description=Run shufflewall on timer

[Timer]
OnUnitActiveSec=60min

[Install]
WantedBy=shufflewall.service
